# launchd/com.ducmn.string-theory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ducmn.string-theory</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ducmn/projects/string-theory/run.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Push to the user's primary (Personal) calendar. The original
             design used a dedicated Tennis calendar shared to a service
             account, but no service-account.json ever existed; we now auth
             as the user via reused OAuth creds (GOOGLE_OAUTH_* below), whose
             calendar.events scope can write to owned calendars but can't
             create new ones. To switch to a dedicated calendar: create one
             in Google Calendar and paste its ID here. -->
        <key>TARGET_CALENDAR_ID</key>
        <string>primary</string>
        <!-- Reuse the foodcycle-slot-scout OAuth desktop-app client + token
             (same Google account, scope calendar.events, has a refresh
             token). calendar_push.build_calendar_service() falls back to
             these when no service account is configured. -->
        <key>GOOGLE_OAUTH_CLIENT_JSON</key>
        <string>/Users/ducmn/projects/foodcycle-slot-scout/credentials.json</string>
        <key>GOOGLE_OAUTH_TOKEN_JSON</key>
        <string>/Users/ducmn/projects/foodcycle-slot-scout/token.json</string>
        <key>BUSY_CALENDAR_IDS</key>
        <string>nguyenmanhduc1309@gmail.com</string>
        <key>BUSY_ICS_URLS</key>
        <string>https://outlook.office365.com/owa/calendar/fa8db5c7e1a047aa8c90cff07e259818@stepstone.com/cd7f0b36024245fa98b8c2d215f2cf932511918829419254462/calendar.ics</string>
        <key>PATH</key>
        <string>/Users/ducmn/Library/Python/3.9/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <!-- Every hour on the hour. Tennis schedules shift through the day
         as preceding matches run long or short, so hourly polling keeps
         the calendar within ~60 min of the upstream schedule rather
         than the old 3h drift window. ~30s per run, ~10 min/day total. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/string-theory.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/string-theory.log</string>
</dict>
</plist>
